# deploy/jamigos.service
# Serviço systemd — para rodar sem Docker, direto com Python.
#
#   sudo cp jamigos.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now jamigos
#   systemctl status jamigos
#
[Unit]
Description=JAMigos - servidor de salas sincronizadas
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jamigos
WorkingDirectory=/opt/jamigos
Environment=PORT=8787
Environment=OJ_DATA=/opt/jamigos/dados/salas.json
# Environment=OJ_TOKEN=suaSenhaAqui
ExecStart=/opt/jamigos/venv/bin/python /opt/jamigos/server/main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
